# systemd/claude-memory-monitor.service
[Unit]
Description=Claude Memory Monitor Daemon
Documentation=https://github.com/yibohub/claude-utilities
After=network.target

[Service]
Type=simple
User=USER_PLACEHOLDER
Group=GROUP_PLACEHOLDER
WorkingDirectory=WORKING_DIR_PLACEHOLDER
Environment="CHECK_INTERVAL=300"
Environment="MEMORY_THRESHOLD=75"
Environment="MAX_CLAUDE_PROCESSES=15"
Environment="AUTO_CLEAN=false"
ExecStart=DAEMON_SCRIPT_PLACEHOLDER
Restart=always
RestartSec=10
StandardOutput=append:LOG_FILE_PLACEHOLDER
StandardError=append:LOG_FILE_PLACEHOLDER

[Install]
WantedBy=multi-user.target
